# deploy/systemd_jtdos_staging.service.example
[Unit]
Description=JTDOS Alpha 0.1 Staging
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/jtdos
EnvironmentFile=/var/www/jtdos/.env.staging
ExecStart=/usr/bin/node backend/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/jtdos

[Install]
WantedBy=multi-user.target
